# files/vpnd.init
#!/bin/sh /etc/rc.common

START=96
EXTRA_COMMANDS="restart reload"
VPND="/bin/vpnd"

start() {
    "$VPND" start
}

stop() {
    "$VPND" stop
}

restart() {
    "$VPND" restart
}

reload() {
    restart
}
